# lib/crontab
# auto gc of local cw data
0 * * * * /home/ubuntu/muad-dweeb/cw/lib/gc.sh
